# deploy/algoforge-web.service
# AlgoForge Next.js Frontend systemd 서비스 파일
#
# 설치 방법:
# sudo cp deploy/algoforge-web.service /etc/systemd/system/
# sudo systemctl daemon-reload
# sudo systemctl enable algoforge-web
# sudo systemctl start algoforge-web
# sudo systemctl status algoforge-web

[Unit]
Description=AlgoForge Next.js Frontend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/algoforge/apps/web
Environment="NODE_ENV=production"
Environment="PORT=3000"
Environment="NEXT_PUBLIC_API_URL=http://algoforge.wonbbo.kro.kr:8080/api"
ExecStart=/usr/bin/pnpm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 리소스 제한 (선택사항)
# MemoryMax=1G
# CPUQuota=100%

[Install]
WantedBy=multi-user.target

