# deploy/algoforge-api.service
# AlgoForge FastAPI Backend systemd 서비스 파일
#
# 설치 방법:
# sudo cp deploy/algoforge-api.service /etc/systemd/system/
# sudo systemctl daemon-reload
# sudo systemctl enable algoforge-api
# sudo systemctl start algoforge-api
# sudo systemctl status algoforge-api

[Unit]
Description=AlgoForge FastAPI Backend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/algoforge
Environment="PATH=/var/www/algoforge/venv/bin"
Environment="PYTHONPATH=/var/www/algoforge"
ExecStart=/var/www/algoforge/venv/bin/uvicorn apps.api.main:app --host 0.0.0.0 --port 8000 --workers 2
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 리소스 제한 (선택사항)
# MemoryMax=2G
# CPUQuota=100%

[Install]
WantedBy=multi-user.target

